# deploy/tg-claude-bot.service
[Unit]
Description=tg-claude Telegram bot (Claude Code interface)
After=network-online.target
Wants=network-online.target

[Service]
User=claude
WorkingDirectory=/opt/tg-claude
ExecStart=/usr/bin/node /opt/tg-claude/dist/bot/main.js
EnvironmentFile=/etc/tg-claude/env
Restart=always
RestartSec=3
# Detached runners must outlive a bot restart:
KillMode=process
# Soft memory cap for the bot + its runners (they share this cgroup): throttle into
# swap past the limit instead of OOM-killing — protects neighbours (VPN) without
# killing a long task. Sized for 961M RAM; raise to ~1300M once the box is at 2G.
MemoryAccounting=yes
MemoryHigh=700M

[Install]
WantedBy=multi-user.target
